# deb/quadlet/unbound.container
[Unit]
Description=Unbound recursive DNS resolver
After=tor-haproxy.service

[Container]
ContainerName=unbound
Image=localhost/unbound:latest
Network=nice-dns.network
IP=172.31.240.251
HealthCmd=dig +short -p 5335 @127.0.0.1 cloudflare.com | grep -Eo '\b[0-9.]{7,15}\b' >/dev/null
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=90s
HealthOnFailure=restart

[Service]
Restart=on-failure
RestartSec=3s
TimeoutStopSec=70

[Install]
WantedBy=default.target
